# services/scx_userland.service
[Unit]
Description=Start scx_userland
ConditionPathExists=/sys/kernel/debug/sched/ext

[Service]
Type=simple
ExecStart=/usr/bin/scx_userland
Restart=always
StandardOutput=null

[Install]
WantedBy=multi-user.target
