# services/scx_nest.service
[Unit]
Description=Start scx_nest
ConditionPathIsDirectory=/sys/kernel/sched_ext

[Service]
Type=simple
ExecStart=scx_nest
Restart=always
StandardError=journal
LogNamespace=sched-ext

[Install]
WantedBy=multi-user.target
