# services/scx_qmap.service
[Unit]
Description=Start scx_qmap
ConditionPathIsDirectory=/sys/kernel/sched_ext

[Service]
Type=simple
ExecStart=scx_qmap
Restart=always
StandardError=journal
LogNamespace=sched-ext

[Install]
WantedBy=multi-user.target
