# services/scx_flatcg.service
[Unit]
Description=Start scx_flatcg
ConditionPathIsDirectory=/sys/kernel/sched_ext

[Service]
Type=simple
ExecStart=scx_flatcg
Restart=always
StandardError=journal
LogNamespace=sched-ext

[Install]
WantedBy=multi-user.target
