# services/scx_rusty.service
[Unit]
Description=Start scx_rusty
ConditionPathExists=/sys/kernel/debug/sched/ext

[Service]
Type=simple
ExecStart=scx_rusty
Restart=always
StandardOutput=null

[Install]
WantedBy=multi-user.target
