# services/scx_nest.service
[Unit]
Description=Start scx_nest
ConditionPathIsDirectory=/sys/kernel/sched_ext
Conflicts=scx_central.service  scx_flatcg.service  scx_pair.service  scx_qmap.service  scx_rustland.service  scx_rusty.service  scx_simple.service  scx_userland.service

[Service]
Type=simple
ExecStart=scx_nest
Restart=always
StandardError=journal
LogNamespace=sched-ext

[Install]
WantedBy=multi-user.target
